# vault/deploy/vaultcare.service
[Unit]
Description=Vaultcare OS Gunicorn
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/vaultcare/backend
Environment="PATH=/var/www/vaultcare/backend/venv/bin"
EnvironmentFile=/var/www/vaultcare/backend/.env
ExecStart=/var/www/vaultcare/backend/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:8000 \
    --timeout 60 \
    vaultcare.wsgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
